# roles/node_backup/templates/node-backup.service.j2
[Unit]
Description=Node backup systemd service

[Service]
Type=oneshot
ExecStart={{ _node_backup_scripts_path }}/common.sh
