# roles/node_backup/templates/node-backup.timer.j2
[Unit]
Description=Node backup systemd timer

[Timer]
{% for time in node_backup_schedule %}
OnCalendar={{ time }}
{% endfor %}
Persistent=true

[Install]
WantedBy=timers.target
